# infra/gcp/viento-backend.service
[Unit]
Description=Backend del Panel Viento (Node.js)
After=network.target

[Service]
Type=simple
User=viento
Group=viento
WorkingDirectory=__APP_DIR__/backend
EnvironmentFile=/etc/viento/backend.env
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
